# install/etc/profile.d/remotepi.sh
export PATH=$PATH:/opt/remotepi/bin:~/.local/bin
