# install/etc/profile.d/remotepi.sh
export PATH=$PATH:/opt/remotepi/bin:~/.local/bin
export QT_QPA_PLATFORMTHEME=gtk2
